# debian/pop-system-updater.postinst
#!/bin/sh

SERVICE="com.system76.SystemUpdater"
SESSION="${SERVICE}.Local"

systemctl daemon-reload

if systemctl is-active "${SERVICE}"; then
    systemctl restart "${SERVICE}"
fi

systemctl enable --now "${SERVICE}"

#DEBHELPER#

exit 0