# debian/pop-system-updater.postinst
#!/bin/sh

SERVICE="com.system76.SystemUpdater"

systemctl daemon-reload

systemctl --root=/ enable "${SERVICE}"

if systemctl is-active "${SERVICE}"; then
    systemctl restart "${SERVICE}"
fi

#DEBHELPER#

exit 0