# data/com.system76.SystemUpdater.Local.service
[Unit]
Description=Local session service for the Pop System Updater
Wants=network-online.target
After=network-online.target

[Service]
Type=dbus
BusName=com.system76.SystemUpdater.Local
ExecStart=/usr/bin/pop-system-updater

[Install]
WantedBy=default.target